# cookbooks/plenv/templates/etc/profile.d/plenv.sh.erb
export PLENV_ROOT=<%= @plenv_root %>
export PATH="$PLENV_ROOT/bin:$PATH"
eval "$(plenv init -)"
