# cookbooks/goenv/templates/etc/profile.d/goenv.sh.erb
export GOENV_ROOT=<%= @goenv_root %>
export PATH="$GOENV_ROOT/bin:$PATH"
eval "$(goenv init -)"
