# cookbooks/goenv/templates/etc/profile.d/goenv.sh.erb
export GOENV_ROOT=<%= @goenv_root %>
export GOPATH=$HOME/.go
export PATH="$GOENV_ROOT/bin:$GOPATH/bin:$PATH"
eval "$(goenv init -)"
