# cookbooks/system/files/etc/cron.hourly/ntpdate
#! /bin/sh
/usr/sbin/ntpdate -u ntp.nict.jp
